# macos/Library/LaunchAgents/com.marc.randomwallpaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.marc.randomwallpaper.plist
  Runs random-wallpaper.sh on login and daily at 9am.
  Missed runs fire when Mac wakes from sleep.

  Deploy:  cd ~/dotfiles && stow macos
  Load:    launchctl load ~/Library/LaunchAgents/com.marc.randomwallpaper.plist
  Reload:  launchctl unload ~/Library/LaunchAgents/com.marc.randomwallpaper.plist && launchctl load ~/Library/LaunchAgents/com.marc.randomwallpaper.plist
  Unload:  launchctl unload ~/Library/LaunchAgents/com.marc.randomwallpaper.plist
  Status:  launchctl list | grep randomwallpaper
  Logs:    cat ~/Library/Logs/random-wallpaper.log
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.marc.randomwallpaper</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/marc.deferranti/.local/bin/random-wallpaper.sh</string>
	</array>
	<!-- Run daily at 9am. Missed runs fire when Mac wakes. -->
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>9</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/marc.deferranti/Library/Logs/random-wallpaper.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/marc.deferranti/Library/Logs/random-wallpaper.log</string>
</dict>
</plist>
